# crontabs
* * * * * /bin/bash /var/pingpong/fetch-urls
