# core/crontabs
* * * * * /bin/sh /scripts/run.sh
# pingpong
