# src/crontabs
* * * * * /bin/sh /scripts/run.sh --follow --append
# pingpong
